# rootfs/etc/inittab
# /etc/inittab
#     init(8) configuration for Debian boot-floppies Root Disk
#
# Copyright (C) 1999,2000 Erik Andersen <andersee@debian.org>
# GPL-2
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# If this is the Linux console, and it's a vga16 frame buffer, 
# force the vga palette to be reset to defaults.
tty1::sysinit:/etc/init.d/check_fb.sh

# main rc script
null::sysinit:/etc/init.d/rcS

# main setup program
::respawn:/sbin/debian-installer

# convenience shell
tty2::askfirst:/bin/sh

#Set up stuff for logging
::sysinit:/bin/dmesg > /dev/tty4
tty3::respawn:/usr/bin/tail -f /var/log/messages
tty4::respawn:/usr/bin/tail -f /proc/kmsg

# Stuff to do before rebooting
::ctrlaltdel:/bin/umount -a -r > /dev/null 2>&1
::ctrlaltdel:/sbin/swapoff -a > /dev/null 2>&1

